# influxdb/crontab
0 */4 * * * /usr/bin/flock -n /tmp/my.lockfile bash /home/import_export.sh
